# debian/gopherd.conffiles
/etc/gopherd/gopherd.conf
/etc/gopherd/gopherdlocal.conf
/etc/init.d/gopherd
/etc/gopherd/startopts
